# webui/systemd/interheart-webui.service
# /opt/interheart/webui/systemd/interheart-webui.service
[Unit]
Description=5echo interheart web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# WebUI kjører som root slik at den kan:
# - lese /var/lib/interheart
# - kjøre interheart CLI uten sudo
User=root
Group=root

WorkingDirectory=/opt/interheart/webui

Environment=WEBUI_BIND=0.0.0.0
Environment=WEBUI_PORT=8088
Environment=PYTHONUNBUFFERED=1

ExecStart=/opt/interheart/webui/.venv/bin/python /opt/interheart/webui/app.py
Restart=always
RestartSec=2

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
