# webui/systemd/interheart-webui.service
[Unit]
Description=5echo interheart web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/interheart/webui

Environment=PORT=8088
Environment=INTERHEART_DB=/var/lib/interheart/state.db
Environment=PYTHONUNBUFFERED=1

ExecStart=/opt/interheart/webui/.venv/bin/python /opt/interheart/webui/app.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
