# webui/systemd/interheart-webui.service
[Unit]
Description=5echo interheart Web UI
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/interheart/webui
Environment=WEBUI_BIND=0.0.0.0
Environment=WEBUI_PORT=8088
ExecStart=/opt/interheart/webui/.venv/bin/python /opt/interheart/webui/app.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
